# wolf-cog-os/payload/etc/systemd/system/cogos-firstboot.service
[Unit]
Description=CoGOS First Boot Governance Setup
DefaultDependencies=no
After=local-fs.target systemd-remount-fs.service
Before=cogos-governance.service

ConditionPathExists=/usr/lib/cogos/firstboot.sh

[Service]
Type=oneshot
ExecStart=/usr/lib/cogos/firstboot.sh
RemainAfterExit=yes
TimeoutStartSec=120
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
