# wolf-cog-os/payload/etc/systemd/system/cogos-first-boot.service
[Unit]
Description=CoGOS post-install first boot takeover
After=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cogos-first-boot
TimeoutStartSec=300

[Install]
WantedBy=graphical.target
