# wolf-cog-os/payload/etc/systemd/system/cogos-boot-hardening.service
[Unit]
Description=CoGOS early boot service hardening (avahi/blueman/udisks)
DefaultDependencies=no
After=local-fs.target
Before=avahi-daemon.service blueman-mechanism.service udisks2.service multi-user.target

ConditionPathExists=/usr/lib/cogos/boot-service-hardening.sh

[Service]
Type=oneshot
ExecStart=/usr/lib/cogos/boot-service-hardening.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
