# wolf-cog-os/payload/etc/systemd/system/cogos-runtime.service
[Unit]
Description=CoGOS governed runtime daemons (native systemd unit)
After=network-online.target local-fs.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cogos-runtime-start
ExecStop=/usr/local/bin/cogos-runtime-stop
RemainAfterExit=yes
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
