# wolf-cog-os/payload/etc/systemd/system/avahi-daemon.service.d/cogos-boot.conf
[Service]
# Prevent multi-minute avahi hangs from blocking multi-user / desktop bring-up.
TimeoutStartSec=45s
TimeoutStopSec=15s
Restart=on-failure
RestartSec=5
