# wolf-cog-os/payload/etc/systemd/system/blueman-mechanism.service.d/cogos-boot.conf
[Unit]
# Bluetooth mechanism needs a graphical session; do not block early boot.
After=graphical.target
Wants=graphical.target

[Service]
TimeoutStartSec=30s
Restart=on-failure
RestartSec=10
