# wolf-cog-os/payload/etc/systemd/system/udisks2.service.d/cogos-boot.conf
[Service]
TimeoutStartSec=60s
Restart=on-failure
RestartSec=5
