# wolf-cog-os/payload/etc/init.d/90cogos
#!/bin/sh
# Legacy SysV wrapper — prefer cogos-runtime.service (native systemd unit).
### BEGIN INIT INFO
# Provides:          cogos
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: CoGOS governed runtime layer
# Description:       Delegates to cogos-runtime-start (native launcher)
### END INIT INFO

case "${1:-start}" in
  start)
    exec /usr/local/bin/cogos-runtime-start
    ;;
  stop)
    exec /usr/local/bin/cogos-runtime-stop
    ;;
  restart)
    /usr/local/bin/cogos-runtime-stop || true
    exec /usr/local/bin/cogos-runtime-start
    ;;
  status)
    if [ -s /opt/cogos/memory/logs/boot_report.json ]; then
      cat /opt/cogos/memory/logs/boot_report.json
    else
      echo "No CoGOS boot report yet. Check /var/log/cogos-service.log"
      exit 1
    fi
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 2
    ;;
esac
